# pkg/DEBIAN/postrm
#!/bin/sh

update-rc.d -f qh3client-daemon remove
