# debian/postrm
#!/bin/sh

service qh3client-daemon stop
update-rc.d -f qh3client-daemon remove
